# pkg/deb/mnt-hotkey-shutdown/DEBIAN/prerm
#!/bin/sh

# Are we being called to remove the package?
if [ "$1" = "remove" ]; then

  # Stop and disable the service
  systemctl disable --now mnt_hotkey_shutdown

fi
